# devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

<<<<<<< HEAD
export CATKIN_TEST_RESULTS_DIR="/afs/ec.auckland.ac.nz/users/z/a/zall747/unixhome/threeohsix/build/test_results"
=======
export CATKIN_TEST_RESULTS_DIR="/afs/ec.auckland.ac.nz/users/c/c/ccha504/unixhome/threeohsix/build/test_results"
>>>>>>> f13782164eafc7b386204ed4d234628a19530598
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
